# ctf-reset.service
[Unit]
Description=Reset Flag CTF Saat Shutdown
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/04-reset.sh

[Install]
WantedBy=shutdown.target reboot.target halt.target
